# coreos-timer/backup-containers.service
[Unit]
Description=Backup all running containers

[Service]
Type=oneshot
ExecStart=/opt/bin/backup-all.sh >> /tmp/backup-all.log'
